# packaging/deb/systemd/isabel.service
[Unit]
Description=Lisa Bluetooth temperature gateway
After=network-online.target bluetooth.service
Wants=network-online.target bluetooth.service

[Service]
Type=simple
EnvironmentFile=/etc/lisa/isabel.env
ExecStart=/usr/bin/isabel
WorkingDirectory=/var/lib/lisa/isabel
Restart=always
RestartSec=5
User=isabel
Group=isabel
AmbientCapabilities=CAP_NET_ADMIN CAP_NET_RAW
CapabilityBoundingSet=CAP_NET_ADMIN CAP_NET_RAW

[Install]
WantedBy=multi-user.target
